# deploy/supervisor/socialdb_vue_flask.conf
[program:socialdb_vue_flask]
command=/usr/local/bin/gunicorn -w10 -b127.0.0.1:8888 api_main:app
directory=/socialdb_vue_flask/server
startsecs=5
stopwaitsecs=0
autostart=true
autorestart=true
